# packaging/systemd/neonstore-backup.timer
[Unit]
Description=Run NeonStore backup nightly
[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
[Install]
WantedBy=timers.target
